# linux/.config/alacritty/themes/autostart/Agente-GAUDI.desktop
[Desktop Entry]
Version=23.0
Name=Agente-GAUDI
Comment=Agente-GAUDI
GenericName=Agente-GAUDI
Exec=/opt/Agente-GAUDI/bin/Agente-GAUDI
Icon=/opt/Agente-GAUDI/lib/Agente-GAUDI.png
Terminal=false
Type=Application
Categories=Network;Application;
StartupNotify=true
MimeType=
X-GNOME-Autostart-enabled=false
