# linux/.config/alacritty/themes/autostart/SCManager.desktop
[Desktop Entry]
Type=Application
Name=SCManager
Exec=/usr/lib/SCMiddleware/SCManager
X-GNOME-Autostart-enabled=false
